# systemd/airflow-webserver.service
[Unit]
Description=Apache Airflow Webserver
After=network.target podft-postgres.service
Wants=podft-postgres.service

[Service]
Type=simple
User=root
Group=root
Environment="AIRFLOW_HOME=/opt/analytics"
Environment="PATH=/opt/analytics/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="AIRFLOW__DATABASE__SQL_ALCHEMY_CONN=postgresql+psycopg2://airflow_user:airflow_pass@localhost:5432/airflow_db"
Environment="AIRFLOW__CORE__EXECUTOR=SequentialExecutor"
Environment="AIRFLOW__CORE__LOAD_EXAMPLES=False"
Environment="AIRFLOW__CORE__DAGS_FOLDER=/opt/analytics/dags"
Environment="AIRFLOW__WEBSERVER__WEB_SERVER_PORT=8080"
Environment="AIRFLOW__WEBSERVER__EXPOSE_CONFIG=True"
ExecStart=/opt/analytics/venv/bin/airflow webserver
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
